# backend/deploy/stock-market.service
# Install: sudo cp stock-market.service /etc/systemd/system/stock-market.service
# Then:  sudo systemctl daemon-reload
#         sudo systemctl enable --now stock-market

[Unit]
Description=Stock Market Backend
Requires=stock-market-env.service
After=network-online.target stock-market-env.service
Wants=network-online.target

[Service]
User=ec2-user
WorkingDirectory=/home/ec2-user/stock-market/backend
EnvironmentFile=/home/ec2-user/stock-market/backend/.env
ExecStart=/home/ec2-user/stock-market/backend/venv/bin/python -u -m uvicorn main:app --host 127.0.0.1 --port 8000 --workers 1 --no-access-log
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
